# etc/crontab
MAILTO="sfloess@localhost"

0 */6 * * * /opt/nas/bin/backups.sh
